# periodic-command-executor.service
[Unit]
Description=Service to read the file periodic-commands.txt

[Service]
Type=simple
ExecStart=/usr/bin/bash /home/sup3r/commands-executor.sh